# deploy/supervisord-executejobs.conf
[program:executejobs]
command=/deploy/start-executejobs.sh
numprocs=1
autostart=true
autorestart=true
